# milo-client/system/milo-client-snapclient.service
[Unit]
Description=Snapclient for Milo Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=milo-client
Group=audio
EnvironmentFile=/var/lib/milo-client/env

ExecStart=/usr/bin/snapclient -h ${MILO_PRINCIPAL_IP} -p 1704 --hostID %H --logsink=system --soundcard camilladsp_input --mixer software

Restart=always
RestartSec=5

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
