# milo-client/system/milo-client-camilladsp.service
[Unit]
Description=Milo Client CamillaDSP Audio Processor
Documentation=https://github.com/HEnquist/camilladsp
After=sound.target
Requires=sound.target

[Service]
Type=simple
User=milo-client
Group=milo-client

# Working directory for config files
WorkingDirectory=/var/lib/milo-client/camilladsp

# Start CamillaDSP with WebSocket server enabled
# -p 1234: WebSocket server port
# -a 127.0.0.1: Listen only on localhost
# -w: Wait for valid config before processing
ExecStart=/usr/local/bin/camilladsp \
    -p 1234 \
    -a 127.0.0.1 \
    -w \
    /var/lib/milo-client/camilladsp/config.yml

# Restart policy
Restart=on-failure
RestartSec=3s

# Resource limits
CPUQuota=80%
MemoryMax=256M

# Security hardening
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=milo-client-camilladsp

[Install]
WantedBy=multi-user.target
